# lesson1/hello-world.service
[Unit]
Description=Hello World Service

[Service]
ExecStart=/usr/bin/hello_world.sh

[Install]
WantedBy=multi-user.target
